# deploy/cloudflared.service
[Unit]
Description=Cloudflare Quick Tunnel (amo2gsheet)
After=network-online.target
Wants=network-online.target
Before=amo2gsheet.service

[Service]
Type=simple
User=root
# Quick tunnel — no account needed. URL is printed to the log on startup.
# Read the assigned URL with: journalctl -u cloudflared | grep trycloudflare
ExecStart=/usr/bin/cloudflared tunnel --no-autoupdate --url http://localhost:8000

Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
TimeoutStartSec=30
TimeoutStopSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

[Install]
WantedBy=multi-user.target
